# etc/supervisord.d/wookteam.ini
[program:wookteam]
directory=/www/wwwroot/team.yunmell.cn
command=php bin/laravels start -i
numprocs=1
autostart=true
autorestart=true
startretries=3
user=root
redirect_stderr=true
stdout_logfile=/www/wwwroot/team.yunmell.cn/%(program_name)s.log